# packaging/build-service/xtreemfs/debian.xtreemfs-server.conffiles
/etc/xos/xtreemfs/dirconfig.properties
/etc/xos/xtreemfs/generate_uuid
/etc/xos/xtreemfs/mrcconfig.properties
/etc/xos/xtreemfs/osdconfig.properties
/etc/init.d/xtreemfs-dir
/etc/init.d/xtreemfs-mrc
/etc/init.d/xtreemfs-osd
