# packaging/build-service/xtreemfs/debian.xtreemfs-client.conffiles
/etc/xos/xtreemfs/default_dir
/etc/init.d/xtreemfs-vivaldi
